# ops/launchd/com.neuropilot.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.neuropilot.watchdog</string>

    <key>WorkingDirectory</key>
    <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>bash</string>
      <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2/ops/launchd/run-watchdog.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>NEUROPILOT_DATA_ROOT</key>
      <string>/Volumes/TradingDrive/NeuroPilotAI</string>
      <key>NEUROPILOT_WATCHDOG_INTERVAL_SEC</key>
      <string>300</string>
      <key>NEUROPILOT_NEARLIVE_LABEL</key>
      <string>com.neuropilot.nearlive</string>
      <key>NEUROPILOT_SMARTLOOP_LABEL</key>
      <string>com.neuropilot.smartloop</string>
      <key>NEUROPILOT_ENV_FILE</key>
      <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2/.env</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>StandardOutPath</key>
    <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2/logs/watchdog.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2/logs/watchdog.err.log</string>
  </dict>
</plist>
