# scripts/com.neuropilot.nightlylab.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.neuropilot.nightlylab</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NEUROPILOT_PROJECT_ROOT</key>
    <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2</string>
    <key>NEUROPILOT_DATA_ROOT</key>
    <string>/Volumes/TradingDrive/NeuroPilotAI</string>
    <key>NEUROPILOT_GDRIVE_BACKUP_ROOT</key>
    <string>/Users/davidmikulis/Library/CloudStorage/GoogleDrive-neuro.pilot.ai@gmail.com/My Drive/Neuro.Pilot.AI/backups</string>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2/engine/scripts/runNightlyResearch.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>20</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2/nightly_lab_launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/davidmikulis/neuro-pilot-ai/neuropilot_trading_v2/nightly_lab_launchd.err.log</string>
</dict>
</plist>
